# crontab
* * * * * /usr/local/bin/php ./artisan schedule:run